# ops/systemd-user/mal-updater-health-check.service
[Unit]
Description=MAL-Updater health-check cycle
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-__MAL_UPDATER_HEALTH_ENV_FILE__
WorkingDirectory=__MAL_UPDATER_REPO_ROOT__
ExecStart=__MAL_UPDATER_REPO_ROOT__/scripts/run_health_check_cycle.sh
