# ops/systemd-user/mal-updater-exact-approved-sync.service
[Unit]
Description=MAL-Updater exact-approved sync cycle
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=<repo-root>
ExecStart=<repo-root>/scripts/run_exact_approved_sync_cycle.sh
