# ops/systemd-user/mal-updater-exact-approved-sync.timer
[Unit]
Description=Run MAL-Updater exact-approved sync about every 8 hours with jitter

[Timer]
Unit=mal-updater-exact-approved-sync.service
OnUnitInactiveSec=6h
RandomizedDelaySec=4h
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
