# ops/systemd-user/mal-updater.service
[Unit]
Description=MAL-Updater background sync daemon
After=default.target

[Service]
Type=simple
WorkingDirectory=__MAL_UPDATER_REPO_ROOT__
Environment=PYTHONPATH=__MAL_UPDATER_REPO_ROOT__/src
EnvironmentFile=-__MAL_UPDATER_SERVICE_ENV_FILE__
ExecStart=python3 -m mal_updater.cli --project-root __MAL_UPDATER_REPO_ROOT__ service-run
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
